# templates/borgbackup.timer.j2
[Unit]
Description=Run daily a borgbackup

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true

[Install]
WantedBy=timers.target
